# biglinux/kde/live-overlay/usr/lib/systemd/system/big-first-boot.service
[Unit]
After=sddm.service

[Service]
ExecStart=/usr/bin/big-first-boot

[Install]
WantedBy=default.target
